# generate-index.timer
[Unit]
Description=Run the service at 5 am

[Timer]
OnCalendar=5:00
Persistent=true

[Install]
WantedBy=timers.target
